# ops/systemd/explorertoken-backend.service
[Unit]
Description=ExplorerToken API
After=network.target redis.service postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=explorer
Group=explorer
WorkingDirectory=/srv/explorertoken/backend
Environment=NODE_ENV=production
EnvironmentFile=/etc/explorertoken/backend.env
ExecStart=/usr/bin/node dist/index.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=20s
# Ensure we log to journald; structured logs emitted by pino/console.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
